# service/wifi-connect-watcher
#!/sbin/openrc-run

# /etc/init.d/wifi_connect_watcher

# Define the script's description
description="Wifi-Connect Watcher for .wifi_connect_mode"

# Path to the watcher script
command="/usr/local/sbin/wifi-connect-watcher.sh"
command_background=true
pidfile="/var/run/wifi-connect-watcher.pid"
name="wifi-connect-watcher"

depend() {
    # Declare dependencies (if any)
    need net
}

start_pre() {
    # Ensure the wifi_connect_watcher.sh is executable
    chmod +x /usr/local/sbin/wifi-connect-watcher.sh
    mkdir -p /var/run/wifi-connect-watcher
    echo "normal_mode" > /var/run/wifi-connect-watcher/request
    echo "normal_mode" > /var/run/wifi-connect-watcher/status
}

start() {
    ebegin "Starting ${name}"
    supervise-daemon ${name} --start --pidfile ${pidfile} --respawn-delay 15 -- ${command}
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    supervise-daemon ${name} --stop --pidfile ${pidfile}
    eend $?
}